# tools/ai-validate/packaging/macos/postinstall
#!/bin/bash
# packaging/macos/postinstall — runs as part of the .pkg install.
#
# pkgbuild stages our bundled .jsx into /tmp/ai-validate-staging/. We
# locate every installed Adobe Illustrator and copy the file into each
# one's Scripts folder so it shows up under File > Scripts. Multi-version
# coexistence (2022 + 2024 side-by-side) is intentionally supported.
#
# Logs to /tmp/ai-validate-postinstall.log AND stderr (which macOS
# captures into /var/log/install.log under the matching pkg-install
# entry) so a silent "no install happened" failure can be diagnosed
# afterwards.

set -euo pipefail

LOG=/tmp/ai-validate-postinstall.log
exec > >(tee -a "$LOG") 2>&1
echo "===== postinstall $(date -u +%FT%TZ) ====="

STAGING="/tmp/ai-validate-staging"
BUNDLE="$STAGING/ai-validate.jsx"

if [[ ! -f "$BUNDLE" ]]; then
    echo "[postinstall] bundle missing at $BUNDLE — aborting"
    exit 1
fi
echo "[postinstall] bundle: $BUNDLE ($(wc -c < "$BUNDLE") bytes)"

shopt -s nullglob
TARGETS=(/Applications/Adobe\ Illustrator*)
if [[ ${#TARGETS[@]} -eq 0 ]]; then
    # No /Applications hit — try mdfind as a fallback (handles
    # custom install locations).
    while IFS= read -r app; do
        TARGETS+=("$app")
    done < <(mdfind 'kMDItemCFBundleIdentifier == "com.adobe.illustrator"' 2>/dev/null)
fi
echo "[postinstall] candidate Illustrator targets: ${TARGETS[*]:-(none)}"

if [[ ${#TARGETS[@]} -eq 0 ]]; then
    osascript -e 'display dialog "No Adobe Illustrator install was detected. ai-validate.jsx is at /tmp/ai-validate-staging/ai-validate.jsx — drop it into Illustrator'\''s Scripts folder manually." buttons {"OK"} default button 1' || true
    exit 0
fi

# macOS Illustrator's Scripts folder lives at <app>/Presets.localized/en_US/Scripts.
# (NOT inside the .app bundle — the directory above /Adobe Illustrator.app
# is the artist-facing root and holds Presets.localized, Plug-ins.localized,
# Cool Extras.localized, etc.) Try the .localized variant first; fall back
# to plain Presets/ for older or stripped builds.
INSTALLED_COUNT=0
for app in "${TARGETS[@]}"; do
    scripts_dir=""
    for candidate in \
        "$app/Presets.localized/en_US/Scripts" \
        "$app/Presets/en_US/Scripts"; do
        if [[ -d "$candidate" ]]; then
            scripts_dir="$candidate"
            break
        fi
    done
    if [[ -z "$scripts_dir" ]]; then
        echo "[postinstall] no Scripts folder under $app — skipped"
        continue
    fi
    cp "$BUNDLE" "$scripts_dir/ai-validate.jsx"
    chmod 644 "$scripts_dir/ai-validate.jsx"
    echo "[postinstall] installed → $scripts_dir/ai-validate.jsx"
    INSTALLED_COUNT=$((INSTALLED_COUNT + 1))
done

if [[ $INSTALLED_COUNT -eq 0 ]]; then
    # Detected Illustrator but couldn't find a Scripts folder under any.
    # Don't fail the install — just surface the issue so the artist knows.
    osascript -e 'display dialog "ai-validate: detected Adobe Illustrator but could not find its Scripts folder. The bundle is at /tmp/ai-validate-staging/ai-validate.jsx — drop it into File > Scripts > Other Script... manually. See /tmp/ai-validate-postinstall.log for details." buttons {"OK"} default button 1' || true
    # Keep staging so the artist can still grab the file.
    exit 0
fi

# Set up the update-check cache + a per-user LaunchAgent that refreshes
# it daily. Why a launchd agent? ExtendScript on macOS Illustrator has
# no silent shell-out — File.execute on .scpt opens Script Editor, on
# .sh opens Terminal — so the script itself can't refresh in
# background. launchd is the proper macOS way: runs in the user's
# session, no UI, fires once per 24h.
#
# Postinstall runs as root for system pkg installs, so all user-facing
# paths (cache, refresh script, LaunchAgent plist) need ownership
# explicitly set to the logged-in user. Pull that user from
# /dev/console's owner — standard macOS idiom.
LOGGED_IN_USER=$(stat -f %Su /dev/console 2>/dev/null || true)
LOGGED_IN_UID=""
if [[ -n "$LOGGED_IN_USER" && "$LOGGED_IN_USER" != "root" ]]; then
    LOGGED_IN_UID=$(id -u "$LOGGED_IN_USER" 2>/dev/null || true)
fi

if [[ -n "$LOGGED_IN_USER" && "$LOGGED_IN_USER" != "root" && -n "$LOGGED_IN_UID" ]]; then
    USER_HOME=$(dscl . -read "/Users/$LOGGED_IN_USER" NFSHomeDirectory 2>/dev/null | awk '{print $2}')
    if [[ -n "$USER_HOME" ]]; then
        CACHE_DIR="$USER_HOME/Library/Application Support/ai-validate"
        REFRESH_SH="$CACHE_DIR/refresh.sh"
        LAUNCH_AGENTS_DIR="$USER_HOME/Library/LaunchAgents"
        PLIST_PATH="$LAUNCH_AGENTS_DIR/com.alexeytuboltsev.ai-validate-updater.plist"

        mkdir -p "$CACHE_DIR" "$LAUNCH_AGENTS_DIR" || true

        # Refresh script — what launchd fires daily (and once at load).
        cat > "$REFRESH_SH" <<'REFRESH_EOF'
#!/bin/bash
# ai-validate manifest refresh — invoked by launchd. Writes the latest
# script-version.json into the local cache so ai-validate.jsx can warn
# the artist about new releases.
set -e
CACHE_DIR="$HOME/Library/Application Support/ai-validate"
mkdir -p "$CACHE_DIR"
curl -fsSL --max-time 10 \
    -o "$CACHE_DIR/script-version-cache.json" \
    "https://raw.githubusercontent.com/AlexeyTuboltsev/house_puzzle/main/tools/ai-validate/script-version.json"
REFRESH_EOF
        chmod 755 "$REFRESH_SH"

        # LaunchAgent plist. RunAtLoad gives an immediate refresh on
        # bootstrap (and on next login); StartInterval=86400 cycles
        # every 24 h while the user session is active. ProcessType:
        # Background hints to launchd this is low-priority.
        cat > "$PLIST_PATH" <<PLIST_EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alexeytuboltsev.ai-validate-updater</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>$REFRESH_SH</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>86400</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/dev/null</string>
    <key>StandardErrorPath</key>
    <string>/dev/null</string>
</dict>
</plist>
PLIST_EOF

        chown -R "$LOGGED_IN_USER" "$CACHE_DIR" "$LAUNCH_AGENTS_DIR" || true
        chmod 644 "$PLIST_PATH"

        # Bootstrap into the user's launchd session. bootout first
        # handles re-installs (bootstrap of an already-loaded plist
        # errors). Both ops tolerate failure — installer should not
        # abort if launchd refuses.
        launchctl bootout "gui/$LOGGED_IN_UID" "$PLIST_PATH" 2>/dev/null || true
        if launchctl bootstrap "gui/$LOGGED_IN_UID" "$PLIST_PATH" 2>>"$LOG"; then
            echo "[postinstall] launchd agent loaded (gui/$LOGGED_IN_UID) — daily refresh active"
        else
            echo "[postinstall] launchd bootstrap failed — agent will pick up at next login"
        fi

        # Backstop: prime the cache synchronously now, in case the
        # bootstrap above failed and RunAtLoad never fires this
        # session.
        if [[ ! -s "$CACHE_DIR/script-version-cache.json" ]]; then
            curl -fsSL --max-time 5 \
                -o "$CACHE_DIR/script-version-cache.json" \
                "https://raw.githubusercontent.com/AlexeyTuboltsev/house_puzzle/main/tools/ai-validate/script-version.json" \
                || echo "[postinstall] cache prime failed (offline?)"
            chown "$LOGGED_IN_USER" "$CACHE_DIR/script-version-cache.json" 2>/dev/null || true
        fi
    fi
fi

# Clean staging on success.
rm -rf "$STAGING"
echo "[postinstall] done — $INSTALLED_COUNT install(s)"
exit 0
